# deploy/muninn-bridge/files/systemd/muninndb.service
[Unit]
Description=MuninnDB - Cognitive Memory Database
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/muninn --daemon --data /root/.muninn/data
ExecStop=/bin/kill -s TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=muninndb

[Install]
WantedBy=multi-user.target
